# hs.cron
# minute hour day-of-month month day-of-week  command
59 * * * * root /opt/hs/.venv/bin/python /opt/hs/scraper.py >> /opt/hs/scraper.log 2>&1
0 * * * * root /usr/bin/node /opt/hs/report.js >> /opt/hs/report.log 2>&1
